# etc/supervisor/supervisord.ini/php-fpm.ini
[program:php]
command=/usr/local/php/sbin/php-fpm -F -c /usr/local/php/lib/php.ini